# sofia/Library/LaunchAgents/com.sofia.indexer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sofia.indexer</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/opt/homebrew/bin/fswatch -o "$HOME/dev/SOFIA" | xargs -n1 -I{} "$HOME/.local/bin/sofia" index --incremental</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/justinramirez/.local/state/sofia/indexer.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/justinramirez/.local/state/sofia/indexer.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
